# rootfs/etc/profile.d/codex-alias.sh
#!/bin/sh
# Convenience alias so additional panes can relaunch Codex quickly
alias codex='NPM_CONFIG_CACHE=/tmp/.npm XDG_CACHE_HOME=/tmp/.cache /usr/local/bin/codex'
