# rootfs/etc/profile.d/ha-cli-aliases.sh
#!/bin/sh
# Convenience helpers for common Home Assistant API and recorder tasks.
alias ha-sqlite='ha-sqlite-ro'
alias ha-db-schema='ha-recorder-schema'
