# ha-codex-terminal/rootfs/etc/profile.d/ha-cli-aliases.sh
#!/bin/sh
# Convenience helpers for common Home Assistant API and recorder tasks.
alias ha-sqlite='ha-sqlite-ro'
alias ha-db-schema='ha-recorder-schema'
alias hass-check-config='hass --script check_config -c /homeassistant'
alias hass-py='ha-python'
